# pingtest/root/usr/bin/pingtest
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
NAME=pingtest
PROG=/usr/bin/pingtest

start_instance() {
	local cfg="$1" param
	config_get_bool restart_wan $cfg restart_wan 1
	config_get domain $cfg domain 'www.baidu.com'
    config_get count $cfg count '5'
    config_get sleep_time $cfg sleep_time '300' 

	procd_open_instance

	procd_set_param command ${PROG} ${domain} ${count} ${sleep_time} ${restart_wan}
	procd_set_param respawn
	procd_close_instance
}

start_service() {
	config_load 'pingtest'
	config_foreach start_instance 'basic'
}

service_triggers() {
      procd_add_reload_trigger "pingtest"
}